# cups/scheduler/org.cups.cupsd.socket
[Unit]
Description=CUPS Scheduler
PartOf=org.cups.cupsd.service

[Socket]
ListenStream=/private/var/run/cupsd

[Install]
WantedBy=sockets.target
