# cups/scheduler/org.cups.cupsd.path
[Unit]
Description=CUPS Scheduler
PartOf=org.cups.cupsd.service

[Path]
PathExists=/private/var/spool/cups/cache/org.cups.cupsd

[Install]
WantedBy=multi-user.target
